# assets/etc/supervisord.conf
[supervisord]
nodaemon=true
user=root
group=root

[program:sshd]
command=/usr/local/sbin/sshd_foreground
autorestart=true

[include]
files = /etc/supervisor/conf.d/*.conf
